# github-webhooks-service/github-webhooks.service
[Unit]
Description=GitHub Webhooks Service
After=network.target

[Service]
User=yoangau
WorkingDirectory=/home/yoangau/cluster/github-webhooks-service
ExecStart=/usr/bin/python github-webhooks-service.py
Restart=always

[Install]
WantedBy=multi-user.target
